# systemd/cli-gateway-sys-executor@.service
[Unit]
Description=CLI Gateway privileged system executor %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Privileged executor must run as root.
User=root
Group=root

WorkingDirectory=/opt/cli-gateway
EnvironmentFile=-/etc/cli-gateway/%i.env
# Prefer project venv, fallback to system Python if missing.
Environment=PATH=/opt/cli-gateway/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env python3 /opt/cli-gateway/system_service_main.py --config /etc/cli-gateway/%i.yaml

Restart=always
RestartSec=5

NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cli-gateway-sys-executor-%i

[Install]
WantedBy=multi-user.target
