# systemd/cli-gateway-system@.service
[Unit]
Description=CLI Gateway (system-admin mode) instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Ops gateway runs in system mode but stays non-root.
User=cli-gateway
Group=cli-gateway

WorkingDirectory=/opt/cli-gateway
EnvironmentFile=-/etc/cli-gateway/%i.env
# Prefer project venv, fallback to system Python if missing.
Environment=PATH=/opt/cli-gateway/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env python3 /opt/cli-gateway/main.py --mode system --instance-id %i --config /etc/cli-gateway/%i.yaml --namespace-paths

Restart=always
RestartSec=5

# Hardening for non-root ops gateway.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/cli-gateway/workspaces /opt/cli-gateway/data /opt/cli-gateway/logs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cli-gateway-system-%i

[Install]
WantedBy=multi-user.target
